# master/debian/postinst
echo "update-rc.d {{package}} defaults || systemctl enable {{package}}"
update-rc.d {{package}} defaults || systemctl enable {{package}}.service
echo "service {{package}} start"
service {{package}} start
