# master/debian/postinst
systemctl enable {{package}}
service {{package}} start
