# master/override/lib/systemd/system/{{package}}.service
[Unit]
Description={{package}} Service

[Service]
Type=simple
ExecStart=/opt/venvs/{{package}}/bin/{{packageService}}
